# etc/profile.d/local-aliases.sh
alias ls='ls -lah --time-style=long-iso'
