# etc/profile.d/local-prompting.sh
export PS1='\n\D{%Y-%m-%d:%H%M %Z}\n[\u@\h \W]\$ '
