# etc/profile.d/local-ports.sh
alias ports='netstat -tulanp | egrep -v "TIME_WAIT|TIME_CLOSE|CLOSE_WAIT|ESTABLISHED"'
