# etc/profile.d/local-backup.sh
backup() { cp $1{,.$(date '+%Y%m%d-%H%M%S')}; }
